# etc/systemd/system/xservices.service
# SyncWatch - Service for SyncWatch

[Unit]
Description=XNAS-xservices (Extended NAS functionality on a linux computer)
After=remote-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/xservices
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
